# scheduler-script/crontab
* * * * * /usr/local/bin/python /app-scheduler/scheduler_by_id.py >> /var/log/cron.log 2>&1
